# orchestration/launchd/com.wolt.assignment.dagster.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wolt.assignment.dagster.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/fhaghighi/Job apply/Wolt/Wolt_Assignment/scripts/dagster_launch_job.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/fhaghighi/Job apply/Wolt/Wolt_Assignment/logs/dagster_daily.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/fhaghighi/Job apply/Wolt/Wolt_Assignment/logs/dagster_daily.err.log</string>
</dict>
</plist>
